# easyprint.supervisor.ini
[program:debug]
directory = /easyprint
command = venv/bin/gunicorn
autostart = false
autorestart = false
startretries = 0