# easyprint.supervisor.ini
[program:easyquest]
directory = /easyprint
command = venv/bin/gunicorn
autostart = true
autorestart = false
startretries = 0