# supervisor.ini
[program:easyquest]
directory = /easyquest
command = venv/bin/gunicorn
autostart = false
autorestart = false
startretries = 0